# etc/systemd/system/lemonbattery.timer
[Unit]
Description=Get battery percentage for lemonbar
Requires=lemonbattery.service

[Timer]
# Time between running
OnBootSec=30
OnUnitActiveSec=1min
Unit=lemonbattery.service

[Install]
WantedBy=timers.target
