# etc/systemd/system/lemonbattery.service
[Unit]
Description=Show battery percentage in lemonbar

[Service]
Type=simple
ExecStart=/home/jessemillar/.dotfiles/scripts/battery
